# hardenedbsd/secadm/pkg-plist
/%%KMODDIR%%/secadm.ko
etc/rc.d/secadm
include/secadm.h
lib/libsecadm.so
lib/libsecadm.so.0
sbin/secadm
